# debian/postinst
#!/bin/sh

set -e 
systemctl daemon-reload
systemctl enable pni-disable-gpu.service
exit 0
