# debian/prerm
#!/bin/sh

set -e 
systemctl disable pni-disable-gpu.service
systemctl daemon-reload
exit 0
