# debian/prerm
#!/bin/bash
set -e
systemctl stop pni-disable-gpu.service || true
systemctl disable pni-disable-gpu.service || true
echo "merhaba" >> /var/cache/prerm-nv
exit 0
